# deploy/metierium-client.service
[Unit]
Description=Metierium Client (Next.js)
After=network.target

[Service]
Type=simple
User=chuck
WorkingDirectory=/home/chuck/projects/metierium/client
ExecStart=/usr/bin/npm start -- -p 3006
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3006

[Install]
WantedBy=multi-user.target
